# telegram-bot/telegram-bot.service
[Unit]
Description=Telegram Bot for Task Management System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/tipa.taska.uz/telegram-bot
Environment="PATH=/var/www/tipa.taska.uz/telegram-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
# ВАЖНО: Загружаем переменные окружения из .env файла
EnvironmentFile=/var/www/tipa.taska.uz/telegram-bot/.env
ExecStart=/var/www/tipa.taska.uz/telegram-bot/venv/bin/python /var/www/tipa.taska.uz/telegram-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
